# hadoop/datanode/conf/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D -o "Port 22" -f /etc/ssh/sshd_config -E /var/log/sshd-server.log
user=root
redirect_stderr=true
stdout_logfile=/var/log/sshd-server.log
autostart=true
autorestart=false
