# yodean/zabbix/zabbix/conf/zabbix.ini
[program:mysql]
command=/usr/libexec/mysqld --user=mysql --basedir=/usr --datadir=/var/lib/mysql --plugin-dir=/usr/lib64/mysql/plugin --log-error=/var/log/mysql/error_mysql.log --pid-file=/var/run/mysql/mysql.pid --socket=/var/lib/mysql/mysql.sock
user=mysql
redirect_stderr=true
stdout_logfile=/var/log/mysql/mysql.log
autostart=true
autorestart=false

[program:apache]
command=/usr/sbin/httpd -D FOREGROUND
user=apache
redirect_stderr=true
stdout_logfile=/var/log/httpd/httpd.log
autostart=true
autorestart=false


