# accountbook/book/conf/book.ini
[program:nginx]
command=/usr/sbin/nginx
user=root
autostart=true
autorestart=true

[program:phpfpm]
command=/usr/sbin/php-fpm -FR
user=root
autostart=true
autorestart=true

; [program:memcached]
; command=/usr/bin/memcached -u memcached -p 11211 -m 256 -c 1024 -I 3m
; user=root
; autostart=true
; autorestart=true

; [program:mysql]
; command=/usr/libexec/mysqld --user=mysql --basedir=/usr --datadir=/var/lib/mysql --plugin-dir=/usr/lib64/mysql/plugin --log-error=/var/log/mariadb/mariadb.log --pid-file=/var/run/mariadb/mariadb.pid --socket=/var/lib/mysql/mysql.sock
; user=mysql
; redirect_stderr=true
; stdout_logfile=/var/log/mariadb.log
; autostart=true
; autorestart=false
